# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable ipmi-exporter.service
systemctl restart ipmi-exporter.service
