# packaging/scripts/preremove.sh
systemctl stop ipmi-exporter.service
systemctl disable ipmi-exporter.service
